# templates/autossh.service.template
[Unit]
Description=AGMind Reverse SSH Tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/autossh -M 0 -N \
  -o "ServerAliveInterval=15" \
  -o "ServerAliveCountMax=3" \
  -o "ExitOnForwardFailure=yes" \
  -o "StrictHostKeyChecking=accept-new" \
  -R 127.0.0.1:__REMOTE_PORT__:localhost:__LOCAL_PORT__ \
  -i __SSH_KEY_PATH__ \
  -p __VPS_SSH_PORT__ tunnel@__VPS_HOST__
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
